# Viking/olad.supervisor.conf
[program:foo]
command=/usr/bin/olad -l 3
autostart=true
autorestart=true
user=alarm
stdout_logfile=/tmp/olad.log
